# scripts/podman/skynet.container.in
# Skynet gateway — Podman Quadlet (rootless)
# Installed by setup-podman.sh into skynet's ~/.config/containers/systemd/
# {{SKYNET_HOME}} is replaced at install time.

[Unit]
Description=Skynet gateway (rootless Podman)

[Container]
Image=skynet:local
ContainerName=skynet
UserNS=keep-id
Volume={{SKYNET_HOME}}/.skynet:/home/node/.skynet
EnvironmentFile={{SKYNET_HOME}}/.skynet/.env
Environment=HOME=/home/node
Environment=TERM=xterm-256color
PublishPort=18789:18789
PublishPort=18790:18790
Pull=never
Exec=node dist/index.js gateway --bind lan --port 18789

[Service]
TimeoutStartSec=300
Restart=on-failure

[Install]
WantedBy=default.target
